# systemd/cronos-docker-restrict.service
[Unit]
Description=Restrict LiteLLM (4001) and dgx-portal (5000) published ports to Traefik + docker bridge + localhost only
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Ports publiés par Docker → filtrés dans la chaîne DOCKER-USER (pas INPUT).
# Sources autorisées : Traefik (172.16.2.7), bridge ai-platform (172.19.0.0/16),
# localhost (127.0.0.1). Tout le reste du LAN est bloqué.
# On insère le DROP puis les ACCEPT (qui passent au-dessus) → ACCEPT prioritaire.
ExecStart=/bin/sh -c '\
for PORT in 4001 5000; do \
  iptables -C DOCKER-USER -p tcp --dport $PORT -j DROP 2>/dev/null || iptables -I DOCKER-USER -p tcp --dport $PORT -j DROP; \
  for SRC in 172.16.2.7 172.19.0.0/16 127.0.0.1; do \
    iptables -C DOCKER-USER -p tcp --dport $PORT -s $SRC -j ACCEPT 2>/dev/null || iptables -I DOCKER-USER -p tcp --dport $PORT -s $SRC -j ACCEPT; \
  done; \
done'
ExecStop=/bin/sh -c '\
for PORT in 4001 5000; do \
  iptables -D DOCKER-USER -p tcp --dport $PORT -j DROP 2>/dev/null || true; \
  for SRC in 172.16.2.7 172.19.0.0/16 127.0.0.1; do \
    iptables -D DOCKER-USER -p tcp --dport $PORT -s $SRC -j ACCEPT 2>/dev/null || true; \
  done; \
done'

[Install]
WantedBy=multi-user.target
